# config/init/upstart/lxc-net.conf.in
description "lxc network"
author "Serge Hallyn <serge.hallyn@canonical.com>"

start on starting lxc
stop on stopped lxc

pre-start exec @LIBEXECDIR@/lxc/lxc-net start
post-stop exec @LIBEXECDIR@/lxc/lxc-net stop
