# config/init/systemd/lxc-monitord.service.in
[Unit]
Description=LXC Container Monitoring Daemon
After=syslog.service network.target
Documentation=man:lxc

[Service]
Type=simple
ExecStart=@LIBEXECDIR@/lxc/lxc-monitord --daemon

[Install]
WantedBy=multi-user.target
